# devproxy.service
[Unit]
Description=devproxy for http proxy server in intranet and direct cases

[Service]
ExecStart=devproxy auto

[Install]
WantedBy=default.target