# static/devproxy.service
[Unit]
Description=devproxy routing directly out or via corporate proxy

[Service]
ExecStart=devproxy auto

[Install]
WantedBy=default.target